# 1.my-config/2-LINUX/2.services/1.headeset-service/99-set-headset-max-volume.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="040b", ATTR{idProduct}=="0897", RUN+="/home/vinicius/.config/systemd/user/scripts/flag-headset-plugged.sh"
