# 1.my-config/2-LINUX/2.services/1.headeset-service/set-headset-max-volume.service
[Unit]
Description=Set headset volume if headset was plugged
After=sound.target

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/scripts/set-headset-max-volume.sh

[Install]
WantedBy=multi-user.target
